# deployment-configs/weather-app.service
[Unit]
Description=Weather Visualization Platform
After=network.target

[Service]
Type=notify
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/weather-visualization
Environment="PATH=/var/www/weather-visualization/venv/bin"
ExecStart=/var/www/weather-visualization/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
